# audit.timer
[UNIT]
Description="timer for the audit.service"

[Timer]

OnCalendar=*-*-* *:00:00
OnCalendar=*-*-* *:30:00

AccuracySec=1s

Unit=audit.service

[Install]
WantedBy=timers.target